# scripts/command-center/master/com.example.cc-master.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Command Center MASTER session (master-loop.sh) — the always-on Opus session
     that is the mobile two-way surface (SECONDARY to the Desktop MCP). Runs on Alpha.
     Copy to ~/Library/LaunchAgents/, replace __KB_ROOT__ / __INSTANCE_NAME__, then:
       launchctl load ~/Library/LaunchAgents/com.example.cc-master.plist
     ⚠️ Deploy AFTER verifying (throwaway) that a headless claude session holds a
     persistent Monitor bus-listener on Alpha — see master-loop.sh's VERIFY-ON-DEPLOY. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.cc-master</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__KB_ROOT__/departments/engineering/command-center/master/master-loop.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>KB_ROOT</key>
        <string>__KB_ROOT__</string>
        <key>CC_MASTER_INSTANCE</key>
        <string>__INSTANCE_NAME__</string>
        <key>CC_MASTER_MODEL</key>
        <string>opus</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>__KB_ROOT__/departments/engineering/command-center</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>/tmp/cc-master.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/cc-master.err</string>
</dict>
</plist>
